# cron/finish_auctions.cron

* * * * * root /usr/bin/python2.7 /path/to/finish_auctions.py
